# scripts/petalia-backup.service
[Unit]
Description=Petalia PostgreSQL Daily Backup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/petalia

# Load environment variables
EnvironmentFile=/home/ec2-user/petalia/.env.production

# Backup script
ExecStart=/home/ec2-user/petalia/scripts/backup-postgres.sh

# Logging
StandardOutput=journal
StandardError=journal

# Timeout (backup can take time for large databases)
TimeoutStartSec=3600

# Restart on failure
OnFailure=petalia-backup-failure@%n.service

[Install]
WantedBy=multi-user.target
