# scripts/petalia-backup-failure@.service
[Unit]
Description=Petalia Backup Failure Notification
PartOf=%i

[Service]
Type=oneshot
User=ec2-user

# Send Slack/email notification on failure
ExecStart=/usr/bin/bash -c ' \
  curl -X POST ${SLACK_WEBHOOK} \
    -H "Content-Type: application/json" \
    -d "{\"text\": \":rotating_light: Petalia backup failed at $(date). Check logs: journalctl -u petalia-backup.service\"}" \
    2>/dev/null || true \
'

StandardOutput=journal
StandardError=journal
